# systemd/user/communal.service
[Unit]
Description=Notification of the data entry

[Service]
Type=simple
ExecStart=/home/ok/.bin/communal
