# systemd/user/lunch.service
[Unit]
Description=A job to test the systemd scheduler

[Service]
Type=simple
ExecStart=%h/python/lunch.py

[Install]
WantedBy=default.target
