# systemd/user/mirrors.service
[Unit]
Description=Update mirrorlist

[Service]
Type=oneshot
ExecStart=/home/ok/git/configs/scripts/mirrors.sh

[Install]
WantedBy=multi-user.target
